# deploy/supervisor_app.conf
[program:cybersecurity_expert]
# Path tempat kode backend Anda berada di VPS
directory=/var/www/cybersecurity-expert-system/backend

# Command untuk menjalankan Gunicorn
command=/var/www/cybersecurity-expert-system/backend/venv/bin/gunicorn app.main:app -c /var/www/cybersecurity-expert-system/backend/deploy/gunicorn_conf.py

# Auto-start & auto-restart
autostart=true
autorestart=true

# Redirect STDERR ke STDOUT
redirect_stderr=true

# Lokasi file log Supervisor
stdout_logfile=/var/log/supervisor/cybersecurity_expert.log
stderr_logfile=/var/log/supervisor/cybersecurity_expert_error.log

# Sesuaikan dengan user VPS yang memiliki akses ke directory (misal: ubuntu/root)
# user=ubuntu

# Set environment variables jika dibutuhkan (opsional, bisa diambil dari .env)
# environment=PATH="/var/www/cybersecurity-expert-system/backend/venv/bin",ENV_FILE="/var/www/cybersecurity-expert-system/backend/.env"
